# 03-kind/kind-dashboard.service
[Unit]
Description=Kind Kubernetes Dashboard Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root
ExecStart=/root/kind-dashboard.sh
Restart=always
RestartSec=10
Environment=KUBECONFIG=/root/.kube/config

[Install]
WantedBy=multi-user.target